# include/etc/profile.d/bin-paths.sh
# shellcheck shell=sh

export PATH="/home/podshell/.local/bin:/home/podshell/.binenv:$PATH"

